# supervisor/motion.conf
[program:motion]
command=/usr/bin/motion
user=tg_motion
autostart=true
autorestart=true
stderr_logfile=/var/log/motion.err.log
stdout_logfile=/var/log/motion.out.log
